# dot_config/autostart/startKiosk.desktop
[Desktop Entry] 

Type=Application

Exec=bash /home/pi/startKiosk
